# production/linux/elements-liquidtestnet.service
[Unit]
Description=Elementsd-liquidtestnet
After=network.target

[Service]
ExecStart=/usr/local/bin/elementsd -daemon -printtoconsole -chain=liquidtestnet -pid=/elements/elements-liquidtestnet.pid
ExecStop=/usr/local/bin/elements-cli stop

Type=forking
PIDFile=/elements/elements-liquidtestnet.pid
Restart=on-failure

User=elements
Group=elements

PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target
