# production/minfee.crontab
@reboot /usr/local/bin/bitcoind >/dev/null 2>&1
